# ansible/roles/restic/templates/restic-backup.service.j2
{{ ansible_managed | comment }}
# restic-backup.service — Restic backup oneshot service
[Unit]
Description=Restic Backup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile={{ restic_env_file }}
Nice={{ restic_nice }}
IOSchedulingClass={{ restic_ionice_class }}

# Backup
ExecStart=/usr/bin/restic backup \
{% for path in restic_backup_paths %}
  {{ path }} \
{% endfor %}
{% for pattern in restic_exclude_patterns %}
  --exclude '{{ pattern }}' \
{% endfor %}
  --verbose

# Forget old snapshots
ExecStartPost=/usr/bin/restic forget \
  --keep-daily {{ restic_keep_daily }} \
  --keep-weekly {{ restic_keep_weekly }} \
  --keep-monthly {{ restic_keep_monthly }} \
  --prune

# Verify repository
ExecStartPost=/usr/bin/restic check

[Install]
WantedBy=multi-user.target
